# etc/crontab
*/15 * * * * /etc/Nuclearite/nuclearite.sh
